# src/p_cam.service
[Unit]
Description=Person Camera Capture Service
After=network.target multi-user.target
# 确保 rkaiq_3A_server (ISP 自动曝光/白平衡) 已启动
# RV1126/RK3568 上服务名通常为 ispserver.service
# 若系统使用其他名称（如 rkaiq_3A.service），请据实修改
After=ispserver.service
Wants=ispserver.service

[Service]
Type=simple
# 延迟 2 秒等待 ISP 3A 完成初始收敛（AE/AWB 稳定）
ExecStartPre=/bin/sleep 2
ExecStart=/home/linaro/p_cam/main --debug
WorkingDirectory=/home/linaro/p_cam
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
